# dist/yukid.service
[Unit]
Description=Sync Local Repos With Remote
After=docker.service
Requires=docker.service
PartOf=docker.serivce

[Service]
User=mirror
ExecStart=/usr/local/bin/yukid

[Install]
WantedBy=multi-user.target
